# crontab
* * * * * * cd /var/www/html && artisan shedule:work >> /dev/null 2>&1
#Empty Line